# nginx/nginx.conf
[program:nginx]
command=/usr/sbin/nginx
environment=
redirect_stderr=true
