# nginx/supervisor.conf
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off; error_log stderr;'
environment=
redirect_stderr=true
